# debian/cron.monthly
#!/bin/sh

[ -x /usr/bin/vbg ] || exit 0

/usr/bin/vbg --quiet
