# debian/cron.monthly
#!/bin/sh

[ -x /usr/bin/vrms ] || exit 0

/usr/bin/vrms --quiet
